# distro/systemd/aios.target
[Unit]
Description=AI-OS.NET Services
Documentation=https://ai-os.net/
Wants=aios-policy-kernel.service aios-evidence-log.service aios-capability-runtime.service aios-sandbox-composer.service aios-vault-broker.service aios-recovery-watchdog.service aios-network-daemon.service aios-ollama.service aios-vllm.service
After=network.target

[Install]
WantedBy=multi-user.target
