# distro/systemd/aios-sgr-daemon.service
[Unit]
Description=AI-OS.NET Service Graph Runtime
Documentation=https://ai-os.net/docs/sgr-daemon
After=network.target aios-capability-runtime.service aios-policy-kernel.service
Requires=aios-capability-runtime.service aios-policy-kernel.service

[Service]
Type=notify
ExecStart=/usr/lib/aios/aios-sgr-daemon --config /etc/aios/config.toml --graph-dir /var/lib/aios/sgr
Restart=always
RestartSec=5
LimitNOFILE=65536
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/aios/sgr /run/aios
NoNewPrivileges=yes
SELinuxContext=system_u:system_r:aios_sgr_daemon_t:s0

[Install]
WantedBy=multi-user.target
