# distro/systemd/aios-network-daemon.service
[Unit]
Description=AI-OS.NET Network Daemon
Documentation=https://ai-os.net/docs/network-daemon
After=network.target aios-capability-runtime.service aios-policy-kernel.service
Requires=aios-capability-runtime.service aios-policy-kernel.service

[Service]
Type=notify
ExecStart=/usr/lib/aios/aios-network-daemon --config /etc/aios/config.toml --state-dir /var/lib/aios/network
Restart=always
RestartSec=5
LimitNOFILE=65536
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/aios/network /run/aios
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_ADMIN CAP_NET_RAW
SELinuxContext=system_u:system_r:aios_network_daemon_t:s0

[Install]
WantedBy=multi-user.target
