# distro/systemd/aios-capability-runtime.service
[Unit]
Description=AI-OS.NET Capability Runtime
Documentation=https://ai-os.net/docs/capability-runtime
After=network.target aios-policy-kernel.service aios-evidence-log.service
Requires=aios-policy-kernel.service aios-evidence-log.service
Wants=aios-vault-broker.service

[Service]
Type=notify
ExecStart=/usr/lib/aios/aios-capability-runtime --config /etc/aios/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
LimitNOFILE=65536
LimitMEMLOCK=infinity
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/aios /run/aios
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SYS_ADMIN
SELinuxContext=system_u:system_r:aios_capability_runtime_t:s0
Environment=AIOS_CONFIG=/etc/aios/config.toml
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
